# scripts/crontab
* * * * * node /home/www/programs/yiwo/task/index.js 1>/home/www/programs/yiwo/log/task.log 2>/home/www/programs/yiwo/log/task.log &