# scripts/crontab
* * * * * node /home/www/bin/MaoDin/task/index.js 1>/home/www/bin/MaoDin/log/task.log 2>/home/www/bin/MaoDin/log/task.log &