# scripts/supervisor.conf
[program:maodin]
command = node index.js 1>log.log 2>log.log &
directory = /home/www/bin/maodin
user = www