# scripts/yiwo.conf
[program:yiwo]
command = node index.js 1>log.log 2>log.log &
directory = /home/www/programs/yiwo
user = www